# deb-structure/DEBIAN/postinst
#!/bin/bash
set -e

sudo systemctl enable thermal_camera.service
sudo systemctl start thermal_camera.service
